# Все/httpd.target
[Unit]
Wants=httpd@conf2.service httpd@conf3.service
