# cpuconcurent.target
[Unit]
Wants=lowprior.service hiprior.service
